# usb.service
[Service]
ExecStart=/home/ubuntu/usbDetectServ
